# resync.service
[Unit]
Description=Resync
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/resync

[Install]
WantedBy=multi-user.target
